# ubuntu_scripts/90-keyboard3.rules
ACTION==”add”, SUBSYSTEM="input", RUN+="/usr/local/bin/xkb_script.sh"

